# tools/gns3/worker-mgmt-failover.service
[Unit]
Description=Worker management backup uplink failover
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/worker-mgmt-failover.sh
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
